# systemd/clear-vim-backups.timer
[Unit]
Description=Timer for the delete vim backups script

[Timer]
OnCalendar=daily

[Install]
WantedBy=default.target
